# util/mma/mma.conf
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Run the mma automation tests"
author          "chromium-os-dev@chromium.org"

start on started boot-services
script
	bash mma_automated_test.sh
end script
